# main.cron
17 1 * * *    python3 /home/de10/kvrk/project/main.py
